# detections/endpoint/Windows/Windows_Autoruns_Run_Key_Value_Set_By_Suspicious Process.yaml
name: Windows Autoruns Run Key Value Set By Suspicious Process
id: Windows_Autoruns_Run_Key_Value_Set_By_Suspicious Process
version: 1
schema_version: "1.0"
date: "2026-03-16"
author: SOCFortress LLC
status: production
type: TTP
description: >
  Detects registry value sets under CurrentVersion\Run or CurrentVersion\Windows
  (Run/Load) keys by suspicious processes including cmd.exe, powershell.exe,
  pwsh.exe, cscript.exe, reg.exe, and regedit.exe. Attackers use these run keys
  to establish persistence, ensuring their payload executes on every user logon.
  Writing to these keys via scripting engines or command-line tools is a strong
  indicator of malicious persistence activity.
  VALIDATED against DESKTOP-AM0HLGE — confirmed reg.exe writing evil.exe path
  to HKU\...\CurrentVersion\Run\EvilPersist.
  MITRE ATT&CK: T1547.001 (Boot or Logon Autostart Execution: Registry Run Keys).

data_source:
  - Sysmon EventID 13 (RegistryEvent Value Set)

search:
  index_pattern: "${INDEX_PATTERN}"
  size: 100
  sort:
    - timestamp: desc
  query:
    bool:
      must:
        - term:
            agent_name: "${AGENT_NAME}"
        - term:
            agent_labels_customer: "${CUSTOMER_CODE}"
        - term:
            data_win_system_eventID: "13"
        - bool:
            should:
              - wildcard:
                  data_win_eventdata_targetObject: "*CurrentVersion\\Run*"
              - bool:
                  must:
                    - wildcard:
                        data_win_eventdata_targetObject: "*CurrentVersion\\Windows*"
                    - bool:
                        should:
                          - wildcard:
                              data_win_eventdata_targetObject: "*\\Run*"
                          - wildcard:
                              data_win_eventdata_targetObject: "*\\Load*"
                        minimum_should_match: 1
            minimum_should_match: 1
        - bool:
            should:
              - wildcard:
                  data_win_eventdata_image: "*\\cmd.exe"
              - wildcard:
                  data_win_eventdata_image: "*\\powershell.exe"
              - wildcard:
                  data_win_eventdata_image: "*\\pwsh.exe"
              - wildcard:
                  data_win_eventdata_image: "*\\cscript.exe"
              - wildcard:
                  data_win_eventdata_image: "*\\reg.exe"
              - wildcard:
                  data_win_eventdata_image: "*\\regedit.exe"
            minimum_should_match: 1
        - range:
            timestamp:
              gte: "${START_TIME}"
              lte: "${END_TIME}"
      must_not:
        - bool:
            should:
              - wildcard:
                  data_win_eventdata_targetObject: "*Discord*"
              - wildcard:
                  data_win_eventdata_targetObject: "*DiscordPTB*"
              - wildcard:
                  data_win_eventdata_targetObject: "*OVPN*"
              - wildcard:
                  data_win_eventdata_targetObject: "*electron*"
              - wildcard:
                  data_win_eventdata_targetObject: "*skype*"
              - wildcard:
                  data_win_eventdata_targetObject: "*ccxprocess*"
              - wildcard:
                  data_win_eventdata_targetObject: "*letsignit*"
              - wildcard:
                  data_win_eventdata_targetObject: "*com.squirrel.Teams.Teams*"
              - wildcard:
                  data_win_eventdata_targetObject: "*ContingencyUninstallNXTSForCurrentUser*"
              - wildcard:
                  data_win_eventdata_targetObject: "*bginfo*"
              - wildcard:
                  data_win_eventdata_targetObject: "*OPENVPN-GUI*"
              - wildcard:
                  data_win_eventdata_targetObject: "*LiveChat*"
              - wildcard:
                  data_win_eventdata_targetObject: "*Mattermost*"
              - wildcard:
                  data_win_eventdata_targetObject: "*Time Doctor*"
              - wildcard:
                  data_win_eventdata_targetObject: "*Flock*"
              - wildcard:
                  data_win_eventdata_targetObject: "*hive*"
              - wildcard:
                  data_win_eventdata_targetObject: "*huddly*"
              - wildcard:
                  data_win_eventdata_targetObject: "*nextiva*"
              - wildcard:
                  data_win_eventdata_targetObject: "*Rambox*"
              - wildcard:
                  data_win_eventdata_targetObject: "*SendAnywhere*"
              - wildcard:
                  data_win_eventdata_targetObject: "*GoogleChromeAutoLaunch*"
              - wildcard:
                  data_win_eventdata_targetObject: "*MicrosoftEdgeAutoLaunch*"
              - wildcard:
                  data_win_eventdata_targetObject: "*TeamsMachineInstaller*"
              - wildcard:
                  data_win_eventdata_targetObject: "*goodaccess*"
              - wildcard:
                  data_win_eventdata_targetObject: "*yammer*"
              - wildcard:
                  data_win_eventdata_targetObject: "*optisigns*"
              - wildcard:
                  data_win_eventdata_targetObject: "*simplifile*"
              - wildcard:
                  data_win_eventdata_targetObject: "*duo*"
              - wildcard:
                  data_win_eventdata_targetObject: "*zalo*"
              - wildcard:
                  data_win_eventdata_targetObject: "*RemoveDuplicateLinks*"
              - wildcard:
                  data_win_eventdata_targetObject: "*krisp*"
              - wildcard:
                  data_win_eventdata_targetObject: "*drata agent*"
              - wildcard:
                  data_win_eventdata_targetObject: "*WDDiscovery*"
              - wildcard:
                  data_win_eventdata_targetObject: "*TimeClockWizard*"
              - wildcard:
                  data_win_eventdata_targetObject: "*SunJavaUPdateSched*"
              - wildcard:
                  data_win_eventdata_targetObject: "*StoreNxDataCreds*"
              - wildcard:
                  data_win_eventdata_targetObject: "*DocsInk*"
              - wildcard:
                  data_win_eventdata_targetObject: "*QtestAutomationHost*"
              - wildcard:
                  data_win_eventdata_targetObject: "*HPSEU_Host_Launcher*"
              - wildcard:
                  data_win_eventdata_targetObject: "*hive*"
              - wildcard:
                  data_win_eventdata_targetObject: "*duo*"
              - wildcard:
                  data_win_eventdata_targetObject: "*zalo*"
              - wildcard:
                  data_win_eventdata_targetObject: "*drata agent*"
            minimum_should_match: 1
        - term:
            data_win_eventdata_details: "(Empty)"
  _source:
    - timestamp
    - agent_name
    - data_win_system_eventID
    - data_win_eventdata_image
    - data_win_eventdata_targetObject
    - data_win_eventdata_details
    - data_win_eventdata_user

parameters:
  AGENT_NAME:
    description: "Target agent/host name"
    type: string
    required: true
    example: "DESKTOP-AM0HLGE"
  CUSTOMER_CODE:
    description: "Customer label for filtering"
    type: string
    required: true
    example: "lab"
  INDEX_PATTERN:
    description: "Elasticsearch index pattern"
    type: string
    required: true
    example: "wazuh-lab_*"
  START_TIME:
    description: "Query start time (ISO 8601 or relative)"
    type: datetime
    required: false
    default: "now-24h"
    example: "2026-03-16 00:00:00.000"
  END_TIME:
    description: "Query end time (ISO 8601 or relative)"
    type: datetime
    required: false
    default: "now"
    example: "2026-03-16 23:59:59.999"

how_to_implement: >
  This detection requires Wazuh agents deployed on Windows endpoints with Sysmon
  installed and configured to log registry value set events (EventID 13).
  Fields data_win_eventdata_image, data_win_eventdata_targetObject, and
  data_win_eventdata_details are confirmed working. VALIDATED — simulation
  confirmed reg.exe writing C:\Windows\Temp\evil.exe to
  HKU\...\CurrentVersion\Run\EvilPersist. Note: Sysmon stores HKCU paths
  as HKU\<SID>\... — the wildcard pattern *CurrentVersion\Run* covers both.

known_false_positives: >
  Legitimate software installers and update mechanisms may write to Run keys
  via cmd.exe or reg.exe. An extensive exclusion list covers known benign
  software (Discord, Teams, Chrome, Edge, OpenVPN, etc.). Additional filters
  may be required for environment-specific software.

references:
  - https://attack.mitre.org/techniques/T1547/001/

response:
  message: >
    Suspicious process wrote to autorun registry key on $agent_name$ by $data_win_eventdata_user$.
    Process: $data_win_eventdata_image$.
    Registry key: $data_win_eventdata_targetObject$.
    Value: $data_win_eventdata_details$.
    MITRE ATT&CK: T1547.001
  risk_score: 60
  severity: high
  risk_objects:
    - field: agent_name
      type: system
      score: 60
    - field: data_win_eventdata_user
      type: user
      score: 60
  threat_objects:
    - field: data_win_eventdata_image
      type: process_name
    - field: data_win_eventdata_targetObject
      type: registry_path

tags:
  analytic_story:
    - Windows Persistence Techniques
    - Windows Registry Abuse
  asset_type: Endpoint
  mitre_attack_id:
    - T1547.001
  custom_tags:
    - run_key
    - persistence
    - registry
    - autoruns
    - defense_evasion
  product:
    - Wazuh
  security_domain: endpoint

graylog:
  query: data_win_system_eventID:"13" AND data_win_eventdata_targetObject:/.*CurrentVersion.*Run.*/ AND (data_win_eventdata_image:/.*\\cmd\.exe/ OR data_win_eventdata_image:/.*\\powershell\.exe/ OR data_win_eventdata_image:/.*\\pwsh\.exe/ OR data_win_eventdata_image:/.*\\cscript\.exe/ OR data_win_eventdata_image:/.*\\reg\.exe/ OR data_win_eventdata_image:/.*\\regedit\.exe/) AND NOT (data_win_eventdata_targetObject:/.*Discord.*/ OR data_win_eventdata_targetObject:/.*DiscordPTB.*/ OR data_win_eventdata_targetObject:/.*OVPN.*/ OR data_win_eventdata_targetObject:/.*electron.*/ OR data_win_eventdata_targetObject:/.*skype.*/ OR data_win_eventdata_targetObject:/.*ccxprocess.*/ OR data_win_eventdata_targetObject:/.*bginfo.*/ OR data_win_eventdata_targetObject:/.*OPENVPN-GUI.*/ OR data_win_eventdata_targetObject:/.*Mattermost.*/ OR data_win_eventdata_targetObject:/.*TeamsMachineInstaller.*/ OR data_win_eventdata_targetObject:/.*GoogleChromeAutoLaunch.*/ OR data_win_eventdata_targetObject:/.*MicrosoftEdgeAutoLaunch.*/ OR data_win_eventdata_targetObject:/.*HPSEU_Host_Launcher.*/ OR data_win_eventdata_targetObject:/.*letsignit.*/ OR data_win_eventdata_targetObject:/.*LiveChat.*/ OR data_win_eventdata_targetObject:/.*Time Doctor.*/ OR data_win_eventdata_targetObject:/.*Flock.*/ OR data_win_eventdata_targetObject:/.*hive.*/ OR data_win_eventdata_targetObject:/.*huddly.*/ OR data_win_eventdata_targetObject:/.*nextiva.*/ OR data_win_eventdata_targetObject:/.*Rambox.*/ OR data_win_eventdata_targetObject:/.*SendAnywhere.*/ OR data_win_eventdata_targetObject:/.*goodaccess.*/ OR data_win_eventdata_targetObject:/.*yammer.*/ OR data_win_eventdata_targetObject:/.*optisigns.*/ OR data_win_eventdata_targetObject:/.*simplifile.*/ OR data_win_eventdata_targetObject:/.*duo.*/ OR data_win_eventdata_targetObject:/.*zalo.*/ OR data_win_eventdata_targetObject:/.*RemoveDuplicateLinks.*/ OR data_win_eventdata_targetObject:/.*krisp.*/ OR data_win_eventdata_targetObject:/.*drata agent.*/ OR data_win_eventdata_targetObject:/.*WDDiscovery.*/ OR data_win_eventdata_targetObject:/.*TimeClockWizard.*/ OR data_win_eventdata_targetObject:/.*SunJavaUPdateSched.*/ OR data_win_eventdata_targetObject:/.*StoreNxDataCreds.*/ OR data_win_eventdata_targetObject:/.*DocsInk.*/ OR data_win_eventdata_targetObject:/.*QtestAutomationHost.*/ OR data_win_eventdata_targetObject:/.*com\.squirrel\.Teams\.Teams.*/ OR data_win_eventdata_targetObject:/.*ContingencyUninstallNXTSForCurrentUser.*/) AND NOT data_win_eventdata_details:"(Empty)"