# detections/endpoint/Powershell/windows_powershell_scheduletask.yaml
name: Windows PowerShell ScheduleTask
id: ddf82fcb-e9ee-40e3-8712-a50b5bf323fc
version: 1
schema_version: "1.0"
date: "2026-03-30"
author: SOCFortress LLC
status: production
type: Anomaly
description: >
  Detects potentially malicious use of PowerShell task scheduling cmdlets
  including New-ScheduledTask, New-ScheduledTaskAction, New-ScheduledTaskSettingsSet,
  New-ScheduledTaskTrigger, Register-ScheduledTask, Set-ScheduledTask,
  Start-ScheduledTask, and Enable-ScheduledTask. Attackers commonly use these
  cmdlets for persistence and remote execution of malicious code. This technique
  is used by Scattered Spider and other threat actors to maintain access, deliver
  additional payloads, or execute ransomware.

data_source:
  - PowerShell Script Block Logging EventID 4104

search:
  index_pattern: "${INDEX_PATTERN}"
  size: 100
  sort:
    - timestamp: desc
  query:
    bool:
      must:
        - term:
            agent_name: "${AGENT_NAME}"
        - term:
            agent_labels_customer: "${CUSTOMER_CODE}"
        - term:
            data_win_system_eventID: "4104"
        - bool:
            should:
              - wildcard:
                  data_win_eventdata_scriptBlockText: "*New-ScheduledTask*"
              - wildcard:
                  data_win_eventdata_scriptBlockText: "*Register-ScheduledTask*"
              - wildcard:
                  data_win_eventdata_scriptBlockText: "*Set-ScheduledTask*"
              - wildcard:
                  data_win_eventdata_scriptBlockText: "*Start-ScheduledTask*"
              - wildcard:
                  data_win_eventdata_scriptBlockText: "*Enable-ScheduledTask*"
              - wildcard:
                  data_win_eventdata_scriptBlockText: "*Register-ClusteredScheduledTask*"
              - wildcard:
                  data_win_eventdata_scriptBlockText: "*Set-ClusteredScheduledTask*"
            minimum_should_match: 1
        - range:
            timestamp:
              gte: "${START_TIME}"
              lte: "${END_TIME}"
  _source:
    - timestamp
    - agent_name
    - data_win_system_eventID
    - data_win_eventdata_scriptBlockText
    - data_win_eventdata_path
    - data_win_eventdata_user

parameters:
  AGENT_NAME:
    description: "Target agent/host name"
    type: string
    required: true
    example: "DESKTOP-AM0HLGE"
  CUSTOMER_CODE:
    description: "Customer label for filtering"
    type: string
    required: true
    example: "lab"
  INDEX_PATTERN:
    description: "Elasticsearch index pattern"
    type: string
    required: true
    example: "wazuh-lab_*"
  START_TIME:
    description: "Query start time (ISO 8601 or relative)"
    type: datetime
    required: false
    default: "now-24h"
    example: "2026-03-30 00:00:00.000"
  END_TIME:
    description: "Query end time (ISO 8601 or relative)"
    type: datetime
    required: false
    default: "now"
    example: "2026-03-30 23:59:59.999"

how_to_implement: >
  Requires PowerShell Script Block Logging enabled on endpoints.
  Enable via Group Policy: Computer Configuration → Administrative Templates
  → Windows Components → Windows PowerShell → Turn on PowerShell Script Block
  Logging.

known_false_positives: >
  Benign administrative tasks such as scheduled maintenance scripts and
  legitimate automation workflows may trigger this detection. Establish a
  baseline of normal scheduled task creation activity and tune accordingly.
  Review the task name, action, and user context to determine legitimacy.

references:
  - https://attack.mitre.org/techniques/T1053/005/
  - https://attack.mitre.org/techniques/T1059/001/
  - https://learn.microsoft.com/en-us/powershell/module/scheduledtasks/
  - https://thedfirreport.com/2023/06/12/a-truly-graceful-wipe-out/

response:
  message: >
    PowerShell scheduled task cmdlet detected on $agent_name$ by
    $data_win_eventdata_user$.
    Script: $data_win_eventdata_scriptBlockText$.
    MITRE ATT&CK: T1053.005, T1059.001
  risk_score: 25
  severity: medium
  risk_objects:
    - field: agent_name
      type: system
      score: 25
    - field: data_win_eventdata_user
      type: user
      score: 25
  threat_objects:
    - field: data_win_eventdata_scriptBlockText
      type: script_block

tags:
  analytic_story:
    - Scheduled Tasks
    - Scattered Spider
  asset_type: Endpoint
  mitre_attack_id:
    - T1053.005
    - T1059.001
  custom_tags:
    - scheduled_task
    - persistence
    - powershell
    - scattered_spider
    - register_scheduledtask
  product:
    - Wazuh
  security_domain: endpoint

graylog:
  query: data_win_system_eventID:"4104" AND (data_win_eventdata_scriptBlockText:/.*New-ScheduledTask.*/ OR data_win_eventdata_scriptBlockText:/.*Register-ScheduledTask.*/ OR data_win_eventdata_scriptBlockText:/.*Set-ScheduledTask.*/ OR data_win_eventdata_scriptBlockText:/.*Start-ScheduledTask.*/ OR data_win_eventdata_scriptBlockText:/.*Enable-ScheduledTask.*/ OR data_win_eventdata_scriptBlockText:/.*Register-ClusteredScheduledTask.*/ OR data_win_eventdata_scriptBlockText:/.*Set-ClusteredScheduledTask.*/)