# detections/endpoint/Linux/linux_service_started_or_enabled.yml
name: Linux Service Started Or Enabled
id: linux-service-started-or-enabled-001
version: 1
schema_version: "1.0"
date: "2026-02-03"
author: SOCFortress LLC
status: production
type: Anomaly

description: >
  Detects the creation or enabling of services on Linux platforms using the systemctl 
  or service tools. Adversaries may create or modify services to maintain persistence 
  or execute malicious payloads on system boot.

data_source:
  - Tetragon
  - Linux Auditd
  - auth.log

search:
  index_pattern: "${INDEX_PATTERN}"
  size: 100
  sort:
    - timestamp: desc
  query:
    bool:
      must:
        - term:
            agent_name: "${AGENT_NAME}"
        - term:
            agent_labels_customer: "${CUSTOMER_CODE}"
        - bool:
            should:
              - wildcard:
                  full_log: "*systemctl*"
              - wildcard:
                  full_log: "*service*"
            minimum_should_match: 1
        - bool:
            should:
              - wildcard:
                  full_log: "* start *"
              - wildcard:
                  full_log: "* enable *"
            minimum_should_match: 1
        - range:
            timestamp:
              gte: "${START_TIME}"
              lte: "${END_TIME}"
      must_not:
        - exists:
            field: data_win_system_eventID
        - wildcard:
            rule_groups: "*windows*"
  _source:
    - timestamp
    - agent_name
    - full_log
    - rule_mitre_id
    - rule_description
    - location

parameters:
  AGENT_NAME:
    description: "Target agent/host name"
    type: string
    required: true
    example: "Amine"
  CUSTOMER_CODE:
    description: "Customer label for filtering"
    type: string
    required: true
    example: "lab"
  START_TIME:
    description: "Query start time (ISO 8601 or relative)"
    type: datetime
    required: false
    default: "now-24h"
    example: "2026-02-03 00:00:00.000"
  END_TIME:
    description: "Query end time (ISO 8601 or relative)"
    type: datetime
    required: false
    default: "now"
    example: "2026-02-03 23:59:59.999"

how_to_implement: >
  This detection requires Wazuh agents deployed on Linux endpoints with Tetragon 
  or auditd process execution logging enabled. Auth.log monitoring provides sudo 
  command visibility. The index pattern should match your Wazuh index naming convention.

known_false_positives: >
  Administrators or network operators may use these commands for legitimate automation 
  and system maintenance. Filter based on known service names, authorized users, and 
  scheduled maintenance windows.

references:
  - https://attack.mitre.org/techniques/T1543/003/

response:
  message: >
    Service start or enable activity detected on $agent_name$. 
    This may indicate persistence establishment via service manipulation. 
    MITRE ATT&CK: $rule_mitre_id$
  risk_score: 42
  severity: medium
  risk_objects:
    - field: agent_name
      type: system
      score: 42
  threat_objects: []

tags:
  analytic_story:
    - Linux Privilege Escalation
    - Linux Persistence Techniques
    - Linux Living Off The Land
    - Scheduled Tasks
  asset_type: Linux
  mitre_attack_id:
    - T1053.006
  product:
    - Wazuh
  security_domain: endpoint
